# root/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

exec cron -f
